# .config/systemd/user/offlineimap.timer
[Unit]
Description=Offlineimap Query Timer


[Timer]
OnCalendar=*:0/15
Unit=offlineimap.service
[Install]
WantedBy=default.target
